# deploy/manifold-companion.service
[Unit]
Description=CAIL OCR Manifold Companion
After=network.target

[Service]
Type=simple
User=smorello.adm@gc.cuny.edu
WorkingDirectory=/data/manifold-companion
ExecStart=/usr/bin/node server/index.js
Environment=NODE_ENV=production
EnvironmentFile=/data/manifold-companion/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
